# vm-aws-creation/aws_linux_windows/roles/canvia.os-monitoring/templates/nmon-monthly.service
[Unit]
After=network.target

[Service]
Type=forking
ExecStart=/scripts/reportes.sh -c -m

[Install]
WantedBy=default.target
