# vm-aws-creation/aws_linux_windows/roles/canvia.os-hardening/templates/compress-history-log.cron
# Comprimir y depurar antiguos logs de historial de comandos
19 0 * * * root /scripts/compress-history-log.sh
